# usr/lib/systemd/system/USBGuardGUIDBus.service
[Unit]
Description=USBGuardGUI DBus Service

[Service]
ExecStart=/usr/bin/USBGuardGUIDBus
Restart=on-failure

[Install]
WantedBy=multi-user.target
