# etc/xdg/autostart/USBGuardGUIApplet.desktop
[Desktop Entry]
Type=Application
Name=USBGuardGUIApplet
Terminal=false
Exec=/usr/bin/USBGuardGUIStartup
